# runcal.supervisor.conf
[program:runcal]
command=python runcal.py
directory=/home/tuxcanfly/build/runcal/
autostart=true
autorestart=true
